# workspaces/rpi-node/rpi-node.service.template
[Unit]
Description=Pegion Defence RPi 4 Node service

[Service]
ExecStart=RPI_NODE_BIN
