# workspaces/rpi-node/rpi-node-api.service.template
[Unit]
Description=Pegion Defence RPi 4 Node API server

[Service]
ExecStart=RPI_NODE_API_BIN
